# overlay/etc/systemd/system/zfs-import-cache.service.d/no-udev-settle.conf
[Unit]
Requires=
After=
After=systemd-udevd.service
