# overlay/etc/systemd/system/docker.service.d/10-pinneos-zfs.conf
[Unit]
# Docker must start after ZFS import so /var/lib/docker is on the ZFS dataset.
# pinneos-zfs-import.service is a oneshot that succeeds even if no pool is found,
# so Docker always starts — it just uses tmpfs if no ZFS pool was imported.
After=pinneos-zfs-import.service
