# overlay/etc/systemd/system/pinneos-usb-mirror-sync.service
[Unit]
Description=PinneOS USB mirror sync (auto-detect mirror disk)
After=multi-user.target

[Service]
Type=oneshot
# Auto-detects any connected non-boot PinneOS USB and syncs both slots.
# Exits cleanly with a log message if no mirror USB is connected.
ExecStart=/usr/lib/homelab/usb-mirror-sync.sh
StandardOutput=journal
StandardError=journal
