# overlay/etc/systemd/system/pinneos-issue.service
[Unit]
Description=PinneOS — generate /etc/issue and /etc/motd with version, hostname and IP
DefaultDependencies=no
# Run after network so hostname -I resolves, but before getty shows the login prompt.
After=NetworkManager-wait-online.service network-online.target
Before=getty.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/homelab/generate-motd.sh

[Install]
WantedBy=getty.target
