# overlay/etc/systemd/system/pinneos-backup-usb-sync@.service
[Unit]
Description=PinneOS backup USB sync for /dev/%i
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/lib/homelab/backup-usb-sync.sh /dev/%i
StandardOutput=journal
StandardError=journal
