# overlay/etc/systemd/system/pinneos-zfs-import.service
[Unit]
Description=PinneOS ZFS pool import
DefaultDependencies=no
After=systemd-udevd.service
# Ensure Docker waits for ZFS so /var/lib/docker comes from ZFS, not tmpfs
Before=docker.service local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/homelab/zfs-import.sh

[Install]
WantedBy=local-fs.target
