# overlay/etc/systemd/system/livecd-alsa-unmuter.service.d/no-udev-settle.conf
[Unit]
After=
After=systemd-udevd.service sound.target sysinit.target
